# assets/systemd/girginospanel-db-backup.service
[Unit]
Description=GirginOSPanel panel DB yedegi

[Service]
Type=oneshot
ExecStart=/usr/local/bin/girginospanel-db-backup
